# systemd/bangumi-mikan-sync.service
[Unit]
Description=Bangumi to Mikan to GuangYa offline-download sync
After=network-online.target guangya-webdav.service
Wants=network-online.target

[Service]
Type=oneshot
ExecStart=/opt/guangya-webdav/venv/bin/python /opt/anime-stack/scripts/bangumi_cloud_download.py --config /opt/anime-stack/config/bangumi-sync.json
WorkingDirectory=/opt/anime-stack
User=root
Group=root
